# docker/cron/crontab
#run python script every minutes
* * * * * root /usr/local/bin/python /app/app/cron.py >> /app/cron.log 2>&1
* * * * * root /usr/local/bin/python /app/app/cron.py retry >> /app/cron.log 2>&1
